# deploy/am-voc-wrapper/am-voc.service.in
[Unit]
Description=AM Vocoder Wrapper service
After=network.target

[Service]
ExecStart=make run-service
WorkingDirectory=${PWD}
StandardOutput=journal
StandardError=journal
Restart=always
User=${USER}

[Install]
WantedBy=multi-user.target